# deployment/sistema_ged.service
[Unit]
Description=Sistema GED - Electronic Document Management System
After=network.target

[Service]
Type=notify
# User and group to run the service
User=www-data
Group=www-data

# Working directory
WorkingDirectory=/path/to/sistema_ged

# Environment file with configuration
EnvironmentFile=/path/to/sistema_ged/.env

# Command to start the service
ExecStart=/path/to/venv/bin/gunicorn -c gunicorn_config.py wsgi:app

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sistema_ged

[Install]
WantedBy=multi-user.target
